# pe.audio.sys/doc/rc.local.example

# switch on your amplifier:
su -l YourUser -c "/home/YourUser/bin/ampli.sh on"


# starts pe.audio.sys
su -l YourUser -c "python3 /home/YourUser/pe.audio.sys/start.py all"


exit 0
