# .install/crontab/auto_update_cronjob

# Queries anacron to update pe.audio.sys from AudioHumLab (master branch) every day
@hourly  /usr/sbin/anacron -s -t $HOME/pe.audio.sys/share/miscel/anacrontab -S $HOME/pe.audio.sys/log
