# .install/crontab/auto_update_cronjob

# Update pe.audio.sys from AudioHumLab (master branch) every day
00  04  *   *   *   sh $HOME/tmp/download_peaudiosys.sh --auto; sh $HOME/tmp/update_peaudiosys.sh --auto
